# config/gunicorn/Changebox.socket
[Unit]
Description=Changebox socket

[Socket]
ListenStream=/run/Changebox.sock

[Install]
WantedBy=sockets.target
